# deploy/systemd/cellclaw-reviewer.service
[Unit]
Description=CellClaw Reviewer Bot
After=network.target

[Service]
Type=simple
User=wangzilu
WorkingDirectory=/path/to/CellClaw
Environment="CELLCRAW_STATE_DIR=/tmp/cellclaw_state"
EnvironmentFile=/path/to/CellClaw/.env.multi
ExecStart=/path/to/CellClaw/venv/bin/python /path/to/CellClaw/main_multi.py reviewer
Restart=always
RestartSec=10
StandardOutput=append:/var/log/cellclaw/reviewer.log
StandardError=append:/var/log/cellclaw/reviewer.error.log

[Install]
WantedBy=multi-user.target
